# server/docker/services.d/grader/run
#!/bin/sh
cd /grader
exec s6-setuidgid grader /usr/local/bin/gunicorn -w 2 grader.wsgi
